# linux-root/etc/NetworkManager/dispatcher.d/30-mtu
#!/bin/sh

if [ "$2" = "up" ]; then
   /usr/sbin/ip link set dev $1 mtu 1390
fi
